# supervisor-app.conf
[program:app-uwsgi]
command = uwsgi /app/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
